# raspios/overlay-rootfs/etc/profile.d/python3.sh
export PATH=/opt/pyenv/versions/3.10.7/bin:$PATH