# raspios/overlay-rootfs/etc/profile.d/python3.sh
export PATH=/opt/pyenv/versions/3.11.2/bin:$PATH
