# raspios/overlay-rootfs/etc/systemd/system/cogip-copilot.service
[Unit]
Description=COGIP Copilot

Wants=network.target
After=syslog.target network-online.target

[Service]
Type=simple
Restart=on-failure
WorkingDirectory=/opt
ExecStart=uv run cogip-copilot
RestartSec=5
EnvironmentFile=/etc/environment

[Install]
WantedBy=multi-user.target
