# raspios/overlay-rootfs/etc/systemd/system/cogip-server.service
[Unit]
Description=COGIP Server

Wants=network.target
After=syslog.target network-online.target

[Service]
Type=simple
Restart=on-failure
WorkingDirectory=/opt
ExecStart=uv run cogip-server
RestartSec=5
EnvironmentFile=/etc/environment

[Install]
WantedBy=multi-user.target
