# raspios/overlay-rootfs/etc/systemd/system/openocd-pami.service
[Unit]
Description=OpenOCD Server for PAMI

Wants=network.target
After=syslog.target network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/openocd -f interface/rpi-cogip-pami.cfg
Restart=on-failure
RestartSec=10
KillMode=process

[Install]
WantedBy=multi-user.target
